# managedl.service
[Unit]
Description=Download manager

[Service]
ExecStart=/home/hendry/go/bin/yyyymmdd

[Install]
WantedBy=default.target
