# managedl.service
[Unit]
Description=Download manager

[Service]
ExecStart=%h/go/bin/yyyymmdd

[Install]
WantedBy=default.target
